# roles/openclaw/templates/claude-max-api-proxy.service.j2
[Unit]
Description=Claude Max API Proxy
After=network.target

[Service]
Type=simple
User={{ username }}
Environment=HOME=/home/{{ username }}
Environment=PATH=/home/{{ username }}/.claude/bin:/usr/local/bin:/usr/bin:/bin
Environment=CLAUDE_CODE_OAUTH_TOKEN={{ claude_oauth_token }}
ExecStart=/usr/local/bin/claude-max-api
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
