# roles/system/templates/docker-image-prune.timer.j2
[Unit]
Description=Run docker image prune on schedule

[Timer]
OnCalendar={{ docker_image_prune_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
