# roles/cup/templates/cup-notify.timer.j2
[Unit]
Description=Run cup-notify on schedule

[Timer]
OnCalendar={{ cup_notify_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
