# roles/immich/templates/immich-storagebox.service.j2
[Unit]
Description=Immich Storage Box Mount
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ username }}
Group={{ username }}
Environment=RCLONE_CACHE_DIR={{ docker_dir }}/immich/rclone-cache
ExecStartPre=/bin/mkdir -p {{ immich_storage_box_mount_path }}
ExecStartPre=/bin/sh -c 'fusermount3 -uz {{ immich_storage_box_mount_path }} || true'
ExecStart=/usr/bin/rclone mount storagebox:{{ immich_storage_box_remote_path }} {{ immich_storage_box_mount_path }} \
  --config {{ docker_dir }}/immich/rclone/rclone.conf \
  --allow-other \
  --cache-dir {{ docker_dir }}/immich/rclone-cache \
  --vfs-cache-mode writes \
  --vfs-cache-max-age 24h \
  --vfs-cache-max-size 2G \
  --vfs-read-chunk-size 32M \
  --vfs-read-chunk-size-limit 256M \
  --dir-cache-time 1h \
  --poll-interval 1m \
  --timeout 1m \
  --contimeout 10s \
  --retries 5 \
  --retries-sleep 2s \
  --low-level-retries 10 \
  --transfers 2 \
  --checkers 2 \
  --sftp-concurrency 1 \
  --sftp-set-modtime=false
ExecStop=/bin/fusermount3 -u {{ immich_storage_box_mount_path }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
