# roles/backup/templates/restic-backup.timer.j2
[Unit]
Description=Run restic backup on schedule

[Timer]
OnCalendar={{ restic_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
